# service/droned/run
#!/bin/sh

/usr/local/bin/droned --config=/etc/drone/drone.toml
